# ulogwatch.timer
[Unit]
Description=Run ulogwatch script

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target